# deploy/server/ping.cron
* * * * * /root/CaCatHead/ping.sh >> /root/cron/log/cron.log 2>&1
